# deploy/systemd/veliorgroup.service
[Unit]
Description=Velior Group Next.js site
After=network.target

[Service]
Type=simple
User=web
Group=web
WorkingDirectory=/home/web/veliorgroup-live
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/npm start --silent
Restart=on-failure
RestartSec=2
StandardOutput=append:/var/log/veliorgroup/app.log
StandardError=append:/var/log/veliorgroup/app.log

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/veliorgroup /home/web/veliorgroup-live
PrivateTmp=true

[Install]
WantedBy=multi-user.target
